# backup/2024-12-18/quadlet/test-scan-project-staging.container
# test-scan-project App - Quadlet Configuration
# Generated by CodeB CLI v2.4.0
# Environment: staging
# Network: codeb-network (DNS-based container communication)

[Unit]
Description=test-scan-project App staging Container
After=network-online.target test-scan-project-staging-postgres.service test-scan-project-staging-redis.service
Requires=test-scan-project-staging-postgres.service test-scan-project-staging-redis.service

[Container]
Image=ghcr.io/${GITHUB_REPOSITORY_OWNER}/test-scan-project:staging
ContainerName=test-scan-project-staging
PublishPort=3003:3000
Network=codeb-network
Environment=NODE_ENV=staging
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
Environment=DATABASE_URL=postgresql://postgres:postgres@test-scan-project-staging-postgres:5432/test_scan_project?schema=public
Environment=REDIS_URL=redis://test-scan-project-staging-redis:6379
EnvironmentFile=/etc/containers/systemd/test-scan-project-staging.env

HealthCmd=curl -sf http://localhost:3000/api/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s
PodmanArgs=--memory=1g --cpus=1

[Service]
Restart=always
RestartSec=10s
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target default.target