# backup/2024-12-18/quadlet/test-scan-project-postgres.container
# test-scan-project PostgreSQL - Quadlet Configuration
# Generated by CodeB CLI v2.4.0
# Environment: production
# Network: codeb-network (DNS-based container communication)

[Unit]
Description=test-scan-project PostgreSQL production Container
After=network-online.target

[Container]
Image=docker.io/library/postgres:15-alpine
ContainerName=test-scan-project-postgres
PublishPort=5432:5432
Network=codeb-network
Environment=POSTGRES_USER=postgres
Environment=POSTGRES_PASSWORD=postgres
Environment=POSTGRES_DB=test_scan_project

Volume=test-scan-project-postgres-data:/var/lib/postgresql/data:Z
HealthCmd=pg_isready -U postgres
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

[Service]
Restart=always
RestartSec=10s
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target default.target