# backup/2024-12-18/quadlet/test-scan-project.container
# test-scan-project App - Quadlet Configuration
# Generated by CodeB CLI v2.4.0
# Environment: production
# Network: codeb-network (DNS-based container communication)

[Unit]
Description=test-scan-project App production Container
After=network-online.target test-scan-project-postgres.service test-scan-project-redis.service
Requires=test-scan-project-postgres.service test-scan-project-redis.service

[Container]
Image=ghcr.io/${GITHUB_REPOSITORY_OWNER}/test-scan-project:latest
ContainerName=test-scan-project
PublishPort=4001:3000
Network=codeb-network
Environment=NODE_ENV=production
Environment=PORT=3000
Environment=HOSTNAME=0.0.0.0
Environment=DATABASE_URL=postgresql://postgres:postgres@test-scan-project-postgres:5432/test_scan_project?schema=public
Environment=REDIS_URL=redis://test-scan-project-redis:6379
EnvironmentFile=/etc/containers/systemd/test-scan-project.env

HealthCmd=curl -sf http://localhost:3000/api/health || exit 1
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s
PodmanArgs=--memory=2g --cpus=2

[Service]
Restart=always
RestartSec=10s
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target default.target