# backup/2024-12-18/quadlet/test-scan-project-redis.container
# test-scan-project Redis - Quadlet Configuration
# Generated by CodeB CLI v2.4.0
# Environment: production
# Network: codeb-network (DNS-based container communication)

[Unit]
Description=test-scan-project Redis production Container
After=network-online.target

[Container]
Image=docker.io/library/redis:7-alpine
ContainerName=test-scan-project-redis
PublishPort=6379:6379
Network=codeb-network

Volume=test-scan-project-redis-data:/data:Z
HealthCmd=redis-cli ping
HealthInterval=30s
HealthTimeout=10s
HealthRetries=3
HealthStartPeriod=60s

[Service]
Restart=always
RestartSec=10s
TimeoutStartSec=300

[Install]
WantedBy=multi-user.target default.target